# examples/kind/supervisord.conf
[supervisord]
nodaemon=true

[program:podman]
command=/usr/bin/podman system service -t 0
stdout_logfile=/var/log/podman.log
stderr_logfile=/var/log/podman_err.log
autostart=true
autorestart=true